# config/network/create-ap.service
[Unit]
Description=Create AP interface & assign static IP

[Service]
Type=oneshot
ExecStart=/usr/bin/createap

[Install]
WantedBy=multi-user.target
